# blaze/launchagents/com.blaze.email-sender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.blaze.email-sender</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts/email_sender.py</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>SUPABASE_URL</key>
    <string>https://briokwdoonawhxisbydy.supabase.co</string>
    <key>SUPABASE_SERVICE_KEY</key>
    <string>eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImJyaW9rd2Rvb25hd2h4aXNieWR5Iiwicm9sZSI6InNlcnZpY2Vfcm9sZSIsImlhdCI6MTc3MTU1Njc2MiwiZXhwIjoyMDg3MTMyNzYyfQ.5V1BsTrqIHGKUUHYJ3PBpL9re_WzKqOzKoQ94dc3me8</string>
  </dict>

  <key>StartInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/_mxappservice/blaze-logs/email-sender.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/_mxappservice/blaze-logs/email-sender.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
