# blaze/launchagents/com.blaze.vector-index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.blaze.vector-index</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts/vector_index.py</string>
    <string>--embed-all</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/_mxappservice</string>
    <key>OPENAI_API_KEY</key>
    <string>sk-proj-***REDACTED***</string>
    <key>SUPABASE_SERVICE_KEY</key>
    <string>eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImJyaW9rd2Rvb25hd2h4aXNieWR5Iiwicm9sZSI6InNlcnZpY2Vfcm9sZSIsImlhdCI6MTc3MTU1Njc2MiwiZXhwIjoyMDg3MTMyNzYyfQ.5V1BsTrqIHGKUUHYJ3PBpL9re_WzKqOzKoQ94dc3me8</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts</string>

  <key>StandardOutPath</key>
  <string>/Users/_mxappservice/blaze-logs/vector_index.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/_mxappservice/blaze-logs/vector_index.err.log</string>

  <!-- Nightly at 9:00 PM — after knowledge-harvest (8 PM) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
