# blaze/launchagents/com.blaze.imsg-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.blaze.imsg-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Developer/CommandLineTools/Library/Frameworks/Python3.framework/Versions/3.9/bin/python3</string>
    <string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts/imsg_watcher.py</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/_mxappservice/logs/imsg_watcher.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/_mxappservice/logs/imsg_watcher.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/_mxappservice</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
