# blaze/launchagents/com.blaze.youtube-learning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>EnvironmentVariables</key>
        <dict>
                <key>HOME</key>
                <string>/Users/_mxappservice</string>
                <key>PATH</key>
                <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        </dict>
        <key>Label</key>
        <string>com.blaze.youtube-learning</string>
        <key>ProgramArguments</key>
        <array>
                <string>/opt/homebrew/bin/python3</string>
                <string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts/youtube_learning_engine.py</string>
        </array>
        <key>StandardErrorPath</key>
        <string>/Users/_mxappservice/blaze-logs/youtube-learning.err.log</string>
        <key>StandardOutPath</key>
        <string>/Users/_mxappservice/blaze-logs/youtube-learning.log</string>
        <key>StartInterval</key>
        <integer>600</integer>
        <key>WorkingDirectory</key>
        <string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts</string>
</dict>
</plist>
