# blaze/launchagents/com.blaze.business-council.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/_mxappservice</string>
		<key>PATH</key>
		<string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>
	<key>Label</key>
	<string>com.blaze.business-council</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts/business_council.py</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/Users/_mxappservice/blaze-logs/business-council.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/_mxappservice/blaze-logs/business-council.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>1</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>/Users/_mxappservice/ACS_CC_AUTOBOT/blaze-v4/ops/scripts</string>
</dict>
</plist>
